# launchd/com.cc-reaper.orphan-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cc-reaper.orphan-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/.cc-reaper/cc-reaper-monitor.sh</string>
    </array>

    <!-- Run every 600 seconds (10 minutes) -->
    <key>StartInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.cc-reaper/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.cc-reaper/logs/launchd-stderr.log</string>
</dict>
</plist>
